# scripts/gitlab-log-extractor.service
[Unit]
Description=GitLab Pipeline Log Extractor Service
After=docker.service
Requires=docker.service

[Service]
Type=forking
WorkingDirectory=/opt/extract-build-logs
ExecStartPre=/bin/sh -c '/usr/bin/docker stop bfa-gitlab-pipeline-extractor 2>/dev/null || true'
ExecStartPre=/bin/sh -c '/usr/bin/docker rm bfa-gitlab-pipeline-extractor 2>/dev/null || true'
ExecStart=/usr/bin/docker run -d \
    --name bfa-gitlab-pipeline-extractor \
    --restart unless-stopped \
    -p 8000:8000 \
    -v /opt/extract-build-logs/logs:/app/logs \
    -v /opt/extract-build-logs/.env:/app/.env:ro \
    gitlab-log-extractor:latest
ExecStop=/usr/bin/docker stop bfa-gitlab-pipeline-extractor
ExecStopPost=/bin/sh -c '/usr/bin/docker rm bfa-gitlab-pipeline-extractor 2>/dev/null || true'
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target
